# yggdrasil-iso/config/includes.chroot/etc/xdg/autostart/yggdrasil-update-check.desktop
[Desktop Entry]
Type=Application
Name=ThorOS Update Check
Comment=Check for ThorOS updates at login and offer them
Exec=/usr/local/bin/yggdrasil-update-check
X-GNOME-Autostart-enabled=true
NoDisplay=true
