# yggdrasil-iso/config/includes.chroot/etc/systemd/system/yggdrasil-firstboot.service
[Unit]
Description=Yggdrasil first-boot setup (GPU detect, model tier, enable app)
After=network-online.target ollama.service
Wants=network-online.target
ConditionPathExists=!/var/lib/yggdrasil/.firstboot-done

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/yggdrasil-firstboot.sh
# No RemainAfterExit / no [Install]: yggdrasil-firstboot.timer drives this and retries every few
# minutes until the model is pulled, so the unit must go inactive after each attempt.
