# yggdrasil-iso/config/includes.chroot/etc/systemd/user/yggdrasil-update-check.service
[Unit]
Description=ThorOS update check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/yggdrasil-update-check
